# linux/roles/conf_backup_server/files/profile.d/custom-aliases.sh
# Manual aliases version 29.01.2024

alias stop='echo STOP > /opt/backup_scripts/backup_from_sftp_to_hdd.status'
alias start='echo START > /opt/backup_scripts/backup_from_sftp_to_hdd.status'
